# package/predixy_1.0.5-1_amd64/lib/systemd/system/predixy.target
[Unit]
Description=Predixy is a high performance and fully featured proxy for redis sentinel and redis clusterredis_auth: false
